# control/postinst
#!/bin/sh
set -e

LOGFILE=/var/log/oticonnect.log
exec >> "$LOGFILE" 2>&1

echo "[oticonnect] postinst start"

# Try to trigger SMHUB backend metadata sync via local API.
RESULT=$(curl -s -w "\n%{http_code}" -X POST \
    "http://localhost:8000/api/smhub/packages/dom.oti.cat/install" 2>/dev/null || true)
HTTP_CODE=$(printf '%s' "$RESULT" | tail -n1)
BODY=$(printf '%s' "$RESULT" | head -n -1)

if [ "$HTTP_CODE" = "200" ]; then
    echo "[oticonnect] SMHUB API registered: $BODY"
else
    echo "[oticonnect] SMHUB API returned $HTTP_CODE: $BODY -- trying direct SQLite insert"
    DBFILE="/opt/smhub-services/config/backend.db"
    if [ -f "$DBFILE" ] && command -v sqlite3 > /dev/null 2>&1; then
        # Remove legacy oticonnect key if present from a previous install
        sqlite3 "$DBFILE" "DELETE FROM apps WHERE key='oticonnect';" 2>/dev/null || true
        sqlite3 "$DBFILE" \
            "INSERT OR REPLACE INTO apps \
               (key, name, description, enabled, start_at_boot, created_at, version, installed, iframe_port, is_service, show_in_sidebar) \
             VALUES \
               ('dom.oti.cat', 'dom.oti.cat', 'Connect SMHUB to oti.cat Home Assistant', \
                1, 1, datetime('now'), '1.0.0-1', 1, 7890, 1, 1);" \
        && echo "[oticonnect] SQLite insert OK" \
        || echo "[oticonnect] SQLite insert failed"
    else
        echo "[oticonnect] sqlite3 not available or DB missing at $DBFILE"
    fi
fi

echo "[oticonnect] Registering in opkg lists so SMHUB backend can discover this package..."
cat > /var/lib/opkg/lists/oticonnect-local << 'EOF'
Package: dom.oti.cat
Version: 1.0.0-1
Architecture: all
Maintainer: oti.cat <hello@oti.cat>
Description: Connect SMHUB Nano MG24 to oti.cat Home Assistant instances
Section: utils
Priority: optional

EOF
echo "[oticonnect] opkg lists entry written"

echo "[oticonnect] Cleaning up legacy stunnel config..."
rm -f /opt/smhub-services/config/stunnel.conf.d/app-oticonnect.conf

echo "[oticonnect] Configuring stunnel proxy (port 27890 -> 7890)..."
STUNNEL_DIR="/opt/smhub-services/config/stunnel.conf.d"
mkdir -p "$STUNNEL_DIR"
cat > "$STUNNEL_DIR/app-dom.oti.cat.conf" << 'EOF'
; Auto-generated by oticonnect postinst. Do not edit manually.
[app-dom.oti.cat]
accept = 0.0.0.0:27890
connect = 127.0.0.1:7890
EOF
kill -HUP "$(cat /var/run/stunnel.pid)" 2>/dev/null && echo "[oticonnect] stunnel reloaded" || echo "[oticonnect] stunnel reload failed"

echo "[oticonnect] Registering with OpenRC and starting service..."
rc-update del oticonnect default 2>/dev/null || true
rc-update add dom.oti.cat default 2>/dev/null && echo "[oticonnect] rc-update OK" || echo "[oticonnect] rc-update failed (may already be added)"
rc-service dom.oti.cat start && echo "[oticonnect] service started" || echo "[oticonnect] service start failed"

echo "[oticonnect] postinst done"
exit 0
